# deploy/systemd/foxboard-cleanup.service
[Unit]
Description=Foxboard import cleanup

[Service]
Type=oneshot
WorkingDirectory=/opt/foxboard
EnvironmentFile=/etc/foxboard.env
ExecStart=/usr/bin/npm run cleanup:imports
User=foxboard
